# waznexserver/misc/ubuntu1204-nginx-gunicorn/etc/init/waznexserver.conf
description "gunicorn for WaznexServer"
start on runlevel [2345]
stop on runlevel [06]
respawn
chdir /opt/WaznexServer
exec /opt/WaznexServer/waznexserver/start_gunicorn.sh
